# build/xanados-installer/configs/system/xanados-cpu-governor.service
# xanadOS Gaming CPU Governor Configuration
# Automatically switches between performance and powersave modes

[Unit]
Description=xanadOS Gaming CPU Governor
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/xanados-cpu-governor performance
ExecStop=/usr/local/bin/xanados-cpu-governor powersave

[Install]
WantedBy=multi-user.target
