# configs/services/xanados-optimizations.service
[Unit]
Description=xanadOS System Optimizations
After=multi-user.target
Before=graphical.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/xanados-system-optimize start
ExecStop=/usr/local/bin/xanados-system-optimize stop

[Install]
WantedBy=multi-user.target
